# deploy/aliyun/flower-position-backup.timer
[Unit]
Description=Daily Flower Position SQLite backup

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
